# perfect_stories_final/techniques/T1053.006/TVR-T1053.006-2026-04-26-007/analytics/evidence_chain.json
{
  "direct_detection": true,
  "correlated_detection": true,
  "evidence_chain": [
    "execution.command_line",
    "telemetry_evidence.key_events[0]",
    "network_telemetry_evidence.key_events[0]",
    "analytic_evidence.sigma[0]",
    "artifact_evidence.files[0]",
    "response_evidence.actions[0]"
  ],
  "anchors": {
    "timestamp_window": {
      "first_seen": "2026-04-18T13:23:14.532315+00:00",
      "last_seen": "2026-04-19T02:32:01.332266+00:00",
      "samples": [
        "2026-04-18T13:23:14.532315+00:00",
        "2026-04-18T14:28:55.729635+00:00",
        "2026-04-18T18:47:50.856560+00:00",
        "2026-04-18T14:35:37.711101+00:00",
        "2026-04-18T18:54:07.947798+00:00",
        "2026-04-22T21:04:03.593735+00:00",
        "2026-04-19T02:32:01.331930+00:00",
        "2026-04-22T21:04:03.595072+00:00"
      ]
    },
    "host_ids": [
      "debian-node-01",
      "mitre-host-f0736972-aea7-4c46-a0a2-144eb5e812c2"
    ],
    "process_ids": [
      4242
    ],
    "ip_addresses": [
      "203.0.113.10"
    ],
    "file_paths": [
      "/etc/systemd/system/timers.target.wants/art-timer.timer",
      "/etc/systemd/system/art-timer.timer",
      "/home/byron/Downloads/Metatron-triune-outbound-gate/atomic-red-team:/opt/atomic-red-team:ro",
      "/home/byron/Downloads/Metatron-triune-outbound-gate/tools/invoke-atomicredteam:/opt/invoke-atomicredteam:ro",
      "/opt/invoke-atomicredteam/Invoke-AtomicRedTeam.psd1",
      "/opt/atomic-red-team/atomics",
      "/tmp/T1053.006.bin",
      "/var/lib/seraph-ai/artifacts/soar_executions_archive.json",
      "/app/backend/data/generated_soar_executions.json",
      "/app/backend/data/soar_executions_archive.json"
    ],
    "hashes": [
      "24ce025939944b498f7600a565287f4d",
      "8493cd51b9efc2cd0ea63a15feb1e15258f761293adc71d134e458aef7bf6d4d",
      "e3b0c44298fc1c149afbf4c8996fb92427ae41e4649b934ca495991b7852b855",
      "3334f4c883554c6388e4c34d931f361b",
      "d4aea76f80dc2b1a2413d3649a8fe898eb1ff40bce1295c03a30163fa7f9bb12",
      "64318071e5e745eb9d6d6b9df4e2aad1",
      "7ebb93abd64f44908eb8dbf638464b8b",
      "aa6712a1ebb7dfbb2f891bd7acd008864a58a1000128de7c722eff4feca99bdf",
      "ab505c23177a4a8ead54b90d54268168"
    ],
    "event_ids": [
      "osq-1776518595-atomic_t1053_006_b1",
      "osq-1776518596-atomic_t1053_006_b2",
      "osq-1776518597-atomic_t1053_006_b3"
    ]
  },
  "story_assessment": {
    "present_layers": [
      "execution",
      "host_telemetry",
      "network_telemetry",
      "detection",
      "artifact",
      "response"
    ],
    "missing_layers": [],
    "layer_count": 6,
    "anchor_counts": {
      "host_ids": 2,
      "process_ids": 1,
      "ip_addresses": 1,
      "file_paths": 10,
      "hashes": 9,
      "event_ids": 3
    },
    "has_anchor_overlap": true,
    "perfect_story": true
  }
}