# metatron_seraph_complete_20260427/integration_evidence/techniques/T1653/TVR-T1653-2026-04-27-002/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1653_CRONTAB",
    "name": "t1653_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  },
  {
    "query_id": "OSQ-T1653_STARTUP",
    "name": "t1653_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  }
]