# perfect_stories_package/techniques/T1053/TVR-T1053-2026-04-26-011/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1053_CRONTAB",
    "name": "t1053_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  },
  {
    "query_id": "OSQ-T1053_AT_JOBS",
    "name": "t1053_at_jobs",
    "query_text": "SELECT path, filename, size, mtime FROM file WHERE path LIKE '/var/spool/cron/atjobs/%' OR path LIKE '/var/spool/at/%';",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  }
]