# metatron_seraph_complete_20260427/integration_evidence/techniques/T1505.001/TVR-T1505.001-2026-04-27-001/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1505_001_CRONTAB",
    "name": "t1505_001_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1505_001_STARTUP",
    "name": "t1505_001_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  }
]