# metatron_seraph_complete_20260427/integration_evidence/techniques/T1543.002/TVR-T1543.002-2026-04-27-001/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1543_002_SYSTEMD_SVC",
    "name": "t1543_002_systemd_svc",
    "query_text": "SELECT path, filename, size, uid, gid, mode, mtime FROM file WHERE path LIKE '/etc/systemd/system/%' OR path LIKE '/usr/lib/systemd/system/%';",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  }
]