# metatron_seraph_complete_20260427/integration_evidence/techniques/T1053.003/TVR-T1053.003-2026-04-27-001/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1053_003_CRONTAB",
    "name": "t1053_003_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  }
]