# metatron_seraph_complete_20260427/integration_evidence/techniques/T1176/TVR-T1176-2026-04-27-004/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1176_CRONTAB",
    "name": "t1176_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1176_STARTUP",
    "name": "t1176_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  }
]