# metatron_full_evidence_bundle_20260424134322/analytics/sigma_matches.json
[
  {
    "detection_basis": "rule_fired_against_osquery_telemetry",
    "level": "medium",
    "matched_event": {
      "action": "added",
      "cmdline": "",
      "image": "",
      "osquery_pack": "crontab_snapshot",
      "timestamp": "Wed Apr 15 08:48:00 2026 UTC",
      "unix_time": "1776242880"
    },
    "rule_file": "file_event_lnx_persistence_cron_files.yml",
    "rule_id": "6c4e2f43-d94d-4ead-b64d-97e53fa2bd05",
    "rule_sha256": "",
    "source": "SigmaHQ/sigma community rules",
    "status": "test",
    "technique": "T1053.003",
    "timestamp": "2026-04-15T08:48:00+00:00",
    "title": "Persistence Via Cron Files"
  },
  {
    "detection_basis": "rule_fired_against_osquery_telemetry",
    "level": "medium",
    "matched_event": {
      "action": "added",
      "cmdline": "",
      "image": "",
      "osquery_pack": "kernel_modules",
      "timestamp": "Wed Apr 15 08:44:15 2026 UTC",
      "unix_time": "1776242655"
    },
    "rule_file": "proc_creation_win_lnx_t1496_hijacking.yml",
    "rule_id": "0l4i5j3k-1g6h-6i9j-0k3l-4h7g8i1j6k7l",
    "rule_sha256": "",
    "source": "SigmaHQ/sigma community rules",
    "status": "experimental",
    "technique": "T1496.001",
    "timestamp": "2026-04-15T08:44:15+00:00",
    "title": "Resource Hijacking via CPU Overload"
  },
  {
    "detection_basis": "rule_fired_against_osquery_telemetry",
    "level": "high",
    "matched_event": {
      "action": "added",
      "cmdline": "",
      "image": "",
      "osquery_pack": "kernel_modules",
      "timestamp": "Wed Apr 15 08:44:15 2026 UTC",
      "unix_time": "1776242655"
    },
    "rule_file": "proc_creation_win_lnx_t1496_hijacking.yml",
    "rule_id": "2n6k7l5m-3i8j-8k1l-2m5n-6j9i0k3l8m9n",
    "rule_sha256": "",
    "source": "SigmaHQ/sigma community rules",
    "status": "experimental",
    "technique": "T1496.003",
    "timestamp": "2026-04-15T08:44:15+00:00",
    "title": "Resource Hijacking via Memory Exhaustion"
  }
]
