# metatron_evidence_bundle_20260427T052729/techniques/T1546.017/TVR-T1546.017-2026-04-27-000/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1546_017_CRONTAB",
    "name": "t1546_017_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  },
  {
    "query_id": "OSQ-T1546_017_STARTUP",
    "name": "t1546_017_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  },
  {
    "query_id": "OSQ-T1546_017_SUID",
    "name": "t1546_017_suid",
    "query_text": "SELECT path, username, groupname, permissions FROM suid_bin;",
    "matched": true,
    "result_count": 15,
    "supporting_event_ids": [
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin"
    ]
  }
]