# metatron_seraph_complete_20260427/integration_evidence/techniques/T1098.005/TVR-T1098.005-2026-04-27-007/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1098_005_CRONTAB",
    "name": "t1098_005_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1098_005_STARTUP",
    "name": "t1098_005_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1098_005_SUID",
    "name": "t1098_005_suid",
    "query_text": "SELECT path, username, groupname, permissions FROM suid_bin;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  }
]