# metatron_evidence_bundle_20260427T052729/techniques/T1546.006/TVR-T1546.006-2026-04-27-000/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1546_006_SUID",
    "name": "t1546_006_suid",
    "query_text": "SELECT path, username, groupname, permissions FROM suid_bin;",
    "matched": true,
    "result_count": 15,
    "supporting_event_ids": [
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin"
    ]
  },
  {
    "query_id": "OSQ-T1546_006_PROC_PRIV",
    "name": "t1546_006_proc_priv",
    "query_text": "SELECT pid, name, path, cmdline, uid, gid FROM processes WHERE uid = 0 AND parent != 0 ORDER BY start_time DESC LIMIT 100;",
    "matched": true,
    "result_count": 1,
    "supporting_event_ids": [
      "osq-1776166996-pack_seraph_processe"
    ]
  },
  {
    "query_id": "OSQ-T1546_006_CRONTAB",
    "name": "t1546_006_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  }
]