# config/osquery.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "disable_logging": "false",
    "log_result_events": "true",
    "schedule_splay_percent": "10",
    "pidfile": "/var/osquery/osquery.pidfile",
    "events_expiry": "3600",
    "database_path": "/var/osquery/osquery.db",
    "verbose": "false",
    "worker_threads": "2",
    "enable_monitor": "true",
    "disable_events": "false",
    "disable_audit": "false",
    "audit_allow_config": "true",
    "host_identifier": "hostname",
    "enable_syslog": "false",
    "audit_allow_sockets": "true",
    "schedule_default_interval": "3600"
  },
  "schedule": {
    "system_info": {
      "query": "SELECT hostname, cpu_brand, physical_memory FROM system_info;",
      "interval": 3600
    },
    "listening_ports": {
      "query": "SELECT pid, port, protocol, family, address FROM listening_ports;",
      "interval": 60
    },
    "process_events": {
      "query": "SELECT pid, parent, path, cmdline FROM processes LIMIT 100;",
      "interval": 30
    },
    "logged_in_users": {
      "query": "SELECT liu.*, p.name, p.cmdline, p.cwd, p.root FROM logged_in_users liu, processes p WHERE liu.pid = p.pid;",
      "interval": 60
    },
    "network_connections": {
      "query": "SELECT p.pid, p.name, p.path, pn.remote_address, pn.remote_port, pn.local_port FROM processes p JOIN process_open_sockets pn ON p.pid = pn.pid WHERE pn.remote_port != 0 LIMIT 50;",
      "interval": 30
    },
    "crontab_snapshot": {
      "query": "SELECT * FROM crontab;",
      "interval": 300
    },
    "startup_items": {
      "query": "SELECT * FROM startup_items;",
      "interval": 3600
    },
    "kernel_modules": {
      "query": "SELECT name, size, used_by FROM kernel_modules;",
      "interval": 300
    },
    "open_files": {
      "query": "SELECT DISTINCT pid, path FROM process_open_files WHERE path NOT LIKE '/proc/%' LIMIT 50;",
      "interval": 60
    },
    "user_accounts": {
      "query": "SELECT username, uid, gid, description, home_directory, shell FROM users WHERE type != 'special' AND shell NOT LIKE '%false' AND shell NOT LIKE '%nologin' AND shell NOT LIKE '%shutdown';",
      "interval": 3600
    },
    "suid_bin": {
      "query": "SELECT * FROM suid_bin;",
      "interval": 3600
    },
    "last_commands": {
      "query": "SELECT uid, command, time FROM shell_history LIMIT 50;",
      "interval": 300
    }
  },
  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT user AS username FROM logged_in_users ORDER BY time DESC LIMIT 1;"
    ]
  },
  "packs": {}
}
