# metatron_seraph_complete_20260427/integration_evidence/techniques/T1547.012/TVR-T1547.012-2026-04-27-000/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1547_012_CRONTAB",
    "name": "t1547_012_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1547_012_STARTUP",
    "name": "t1547_012_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1547_012_SUID",
    "name": "t1547_012_suid",
    "query_text": "SELECT path, username, groupname, permissions FROM suid_bin;",
    "matched": true,
    "result_count": 15,
    "supporting_event_ids": [
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin",
      "osq-1776244466-suid_bin"
    ],
    "osquery_evidence_type": "direct_match"
  }
]