# metatron_seraph_complete_20260427/integration_evidence/techniques/T1547.001/TVR-T1547.001-2026-04-27-006/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1547_001_STARTUP",
    "name": "t1547_001_startup",
    "query_text": "SELECT name, path, source, status, type FROM startup_items;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": []
  }
]