# metatron_seraph_complete_20260427/integration_evidence/techniques/T1546.013/TVR-T1546.013-2026-04-27-001/analytics/osquery_correlations.json
[
  {
    "query_id": "OSQ-T1546_013_SUID",
    "name": "t1546_013_suid",
    "query_text": "SELECT path, username, groupname, permissions FROM suid_bin;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1546_013_PROC_PRIV",
    "name": "t1546_013_proc_priv",
    "query_text": "SELECT pid, name, path, cmdline, uid, gid FROM processes WHERE uid = 0 AND parent != 0 ORDER BY start_time DESC LIMIT 100;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  },
  {
    "query_id": "OSQ-T1546_013_CRONTAB",
    "name": "t1546_013_crontab",
    "query_text": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
    "matched": false,
    "result_count": 0,
    "supporting_event_ids": [],
    "osquery_evidence_type": "mapped_query_only"
  }
]