# perfect_stories_final/atomic_runs/run_8cc5504a557d4ad7a511729b0f730c83.json
{
  "run_id": "8cc5504a557d4ad7a511729b0f730c83",
  "job_id": "fix-remediation-t1053-003",
  "job_name": "T1053.003-remediation-cron",
  "status": "success",
  "message": "Technique executed successfully",
  "techniques": [
    "T1053.003"
  ],
  "detected_techniques": [
    "T1053.003"
  ],
  "techniques_executed": [
    "T1053.003"
  ],
  "command": "Invoke-AtomicTest T1053.003",
  "exit_code": 0,
  "stdout": "PathToAtomicsFolder = /opt/atomic-red-team/atomics\n\nExecuting test: T1053.003-2 Cron - Add script to all cron subfolders\nbash: line 1: /etc/cron.hourly/persistevil: No such file or directory\nbash: line 1: /etc/cron.monthly/persistevil: No such file or directory\nbash: line 1: /etc/cron.weekly/persistevil: No such file or directory\nExit code: 1\nDone executing test: T1053.003-2 Cron - Add script to all cron subfolders\nExecuting test: T1053.003-3 Cron - Add script to /etc/cron.d folder\nExit code: 0\nDone executing test: T1053.003-3 Cron - Add script to /etc/cron.d folder\nExecuting test: T1053.003-4 Cron - Add script to /var/spool/cron/crontabs/ folder\nExit code: 0\nDone executing test: T1053.003-4 Cron - Add script to /var/spool/cron/crontabs/ folder\n",
  "stderr": "",
  "started_at": "2026-04-18T11:31:59.800104Z",
  "finished_at": "2026-04-18T11:31:59.800104Z",
  "dry_run": false,
  "real_execution": true,
  "execution_mode": "Exec",
  "sandbox": "docker-network-none-cap-drop-all",
  "atomics_path": "/opt/atomic-red-team/atomics",
  "remediation_run": true
}