# metatron_seraph_complete_20260427/integration_evidence/techniques/T1543.002/TVR-T1543.002-2026-04-27-002/analytics/evidence_chain.json
{
  "direct_detection": true,
  "correlated_detection": true,
  "evidence_chain": [
    "execution.command_line",
    "telemetry_evidence.key_events[0]",
    "network_telemetry_evidence.key_events[0]",
    "analytic_evidence.sigma[0]",
    "artifact_evidence.files[0]",
    "response_evidence.actions[0]"
  ],
  "anchors": {
    "timestamp_window": {
      "first_seen": "2026-04-25T07:49:21.433086+00:00",
      "last_seen": "1777117750.774221",
      "samples": [
        "2026-04-25T07:49:21.433086+00:00",
        "2026-04-18T14:39:35.526616+00:00",
        "2026-04-18T19:09:27.740772+00:00",
        "2026-04-18T13:24:00.091396+00:00",
        "2026-04-18T14:31:57.153475+00:00",
        "2026-04-18T18:50:44.287828+00:00",
        "2026-04-18T18:57:02.036865+00:00",
        "2026-04-22T21:13:39.047811+00:00"
      ]
    },
    "host_ids": [
      "debian-node-01",
      "mitre-host-f0736972-aea7-4c46-a0a2-144eb5e812c2"
    ],
    "process_ids": [
      4242
    ],
    "ip_addresses": [
      "203.0.113.10",
      "192.168.101.183",
      "192.168.101.1",
      "142.251.47.99"
    ],
    "file_paths": [
      "/opt/atomic-red-team/atomics",
      "/7",
      "/etc/systemd/system/default.target.wants/art-systemd-service.service",
      "/etc/systemd/system/art-systemd-service.service",
      "/usr/local/etc/rc.d/art-test:",
      "/usr/local/etc",
      "/usr/local/etc/rc.d/art-test",
      "/etc/init.d/art-test:",
      "/usr",
      "/usr/lib/systemd/systemd-sysv-i",
      "/usr/lib/systemd/systemd-sysv-install",
      "/usr/lib/systemd/systemd-sysv-ins"
    ],
    "hashes": [
      "eeacc3100b8d4c478f140eddd0b1ab84",
      "6ccbbf28fc34df63e7ce96bf9231575e1f82d436ec874eb9b8f1a967afeebd89",
      "e3b0c44298fc1c149afbf4c8996fb92427ae41e4649b934ca495991b7852b855"
    ],
    "event_ids": [
      "osq-1777103362-atomic_t1543_002_b1",
      "osq-1776523176-atomic_t1543_002_b1",
      "osq-1776523177-atomic_t1543_002_b2",
      "osq-1776523178-atomic_t1543_002_b3",
      "osq-1776539368-atomic_t1543_002_b1",
      "osq-1776539369-atomic_t1543_002_b2",
      "osq-1776539370-atomic_t1543_002_b3",
      "osq-1776518641-atomic_t1543_002_b1",
      "osq-1776518642-atomic_t1543_002_b2",
      "osq-1776518643-atomic_t1543_002_b3",
      "osq-1776522718-atomic_t1543_002_b1",
      "osq-1776522719-atomic_t1543_002_b2",
      "osq-1776522720-atomic_t1543_002_b3",
      "osq-1776538245-atomic_t1543_002_b1",
      "osq-1776538246-atomic_t1543_002_b2",
      "osq-1776538247-atomic_t1543_002_b3",
      "osq-1776538623-atomic_t1543_002_b1",
      "osq-1776538624-atomic_t1543_002_b2",
      "osq-1776538625-atomic_t1543_002_b3"
    ]
  },
  "story_assessment": {
    "present_layers": [
      "execution",
      "host_telemetry",
      "network_telemetry",
      "detection",
      "artifact",
      "response"
    ],
    "missing_layers": [],
    "layer_count": 6,
    "anchor_counts": {
      "host_ids": 2,
      "process_ids": 1,
      "ip_addresses": 4,
      "file_paths": 12,
      "hashes": 3,
      "event_ids": 19
    },
    "has_anchor_overlap": true,
    "perfect_story": true
  }
}